# docker/php-cli/supervisor.d/crond.ini
[program:crond]
command=/usr/sbin/crond -f -l 8
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=5MB