# docker/nginx/supervisor.d/nginx.ini
[program:nginx]
command=nginx -g 'daemon off;'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=5MB
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=5MB